# build/supervisor/supervisord.conf
[supervisord]
loglevel=debug  # Increase the loglevel to debug for better insights
logfile=/dev/stdout  # Direct Supervisor logs to stdout
stderr_logfile=/dev/stderr
nodaemon=true
pidfile=/var/run/supervisord.pid
user=root

[group:laravel-worker]
priority=999
programs=nginx,php8-fpm

[program:nginx]
priority=10
autostart=true
autorestart=true
startsecs=5  # Ensure nginx has time to start properly
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stderr_events_enabled=true
command=/usr/sbin/nginx -g 'daemon off;'

[program:php8-fpm]
priority=5
autostart=true
autorestart=true
startsecs=5  # Ensure php-fpm has time to start properly
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
command=/usr/local/sbin/php-fpm -R
